# iptables_lxc.path
[Unit]
Description=Container iptables ip update path

[Path]
PathModified=/tmp/ip_file.sh
Unit=iptables_lxc.service

[Install]
WantedBy=multi-user.target
