# iptables_arm.path
[Unit]
Description=Container iptables ip update path

[Path]
PathModified=/tmp/ip_file.sh
Unit=iptables_arm.service

[Install]
WantedBy=multi-user.target