# lxc.path
[Unit]
Description=Container Path Based initialization

[Path]
PathModified=/tmp/.lxcenabled
Unit=lxc.service

[Install]
WantedBy=multi-user.target
